# src/emulator/image/basic-boot
#!/sbin/openrc-run

description="Run Essential Boot Scripts"

start() {
  ebegin "Running Essential Boot Scripts"
  mount / -o remount,rw
  mount /puter
  eend $?
}

stop() {
  ebegin "Stopping Essential Boot Scripts"
  eend $?
}
