# src/emulator/image/twisp-service
#!/sbin/openrc-run

description="twisp daemon"
command="/bin/twisp"
command_args="--pty /dev/hvc0"
pidfile="/var/run/twisp.pid"
command_background="yes"
start_stop_daemon_args="--background --make-pidfile"

depend() {
    need localmount
    after bootmisc
}

start() {
    ebegin "Starting ${description}"
    start-stop-daemon --start --pidfile "${pidfile}" --background --exec ${command} -- ${command_args}
    eend $?
}

stop() {
    ebegin "Stopping ${description}"
    start-stop-daemon --stop --pidfile "${pidfile}"
    eend $?
}
